# supervisor/avahi-server-browse.conf
[program:avahi-server-browse]
command=/home/flux/samson/avahi-server-browse/server.py
autostart=true
autorestart=true
startretries=3
user=flux
